# deploy/essenciel-ingestion.timer
[Unit]
Description=Ingestion carburants deux fois par jour (minuit et midi)
Requires=essenciel-ingestion.service

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target
